# autostart/compton.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/compton