# autostart/fuckingBell.desktop
[Desktop Entry]

Type=Application
Exec=xset b off
