# autostart/pulseaudio.desktop
[Desktop Entry]

Type=Application
Exec=pulseaudio --start²:wq